# rootfs/etc/s6-overlay/s6-rc.d/slskd/run
#!/usr/bin/with-contenv sh

# run with limited permissions
exec s6-setuidgid $PUID:$PGID /app/bin/slskd --app-dir=/config --config=/config/slskd.yml
